# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
username=devpass              ; (default is no username (open server))
password={SHA}dbecba2b0bfde3de830fa6b9fbddf2717f5a8f4e             ; (default is no password (open server))

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=1MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
username=devpass              ; should be same as http_username if set
password={SHA}dbecba2b0bfde3de830fa6b9fbddf2717f5a8f4e                ; should be same as http_password if set
stopwaitsecs=10
priority=1000


[program:global_counter]
directory=/home/ubuntu/RedisGlobalKeyPOC
command=python3 main.py
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/home/ubuntu/stdout.log
stderr_logfile=/home/ubuntu/stderr.log
autostart=true
autorestart=true
startsecs=10
numprocs=2
